# supervisor-app.conf
[program:uwsgi]
command=uwsgi --emperor /var/www/mspray/uwsgi.ini --master
autostart=true
autorestart=true
stdout_logfile=/var/log/mspray.log
redirect_stderr=True
